# fs-skel/etc/local.d/aidos-setup.start
#!/bin/sh
# First-boot user account setup
/opt/aidos/scripts/firstboot.sh
